# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
user=user

[program:ai-engine]
command=uvicorn main:app --host 127.0.0.1 --port 8000
directory=/home/user/app/ai-engine
autostart=true
autorestart=true
stdout_logfile=/tmp/ai-engine.log
stderr_logfile=/tmp/ai-engine.err
environment=PATH="/usr/local/bin:%(ENV_PATH)s"

[program:backend]
command=node server.js
directory=/home/user/app/backend
autostart=true
autorestart=true
stdout_logfile=/tmp/backend.log
stderr_logfile=/tmp/backend.err
environment=PORT="5000",NODE_ENV="production",JWT_SECRET="cybershield_hf_secret",AI_ENGINE_URL="http://127.0.0.1:8000"

[program:frontend]
command=node_modules/.bin/next start -p 3000
directory=/home/user/app/frontend
autostart=true
autorestart=true
stdout_logfile=/tmp/frontend.log
stderr_logfile=/tmp/frontend.err
environment=NEXT_TELEMETRY_DISABLED="1",NEXT_PUBLIC_API_URL="/api"

[program:nginx]
command=nginx -c /etc/nginx/nginx.conf
autostart=true
autorestart=true
stdout_logfile=/tmp/nginx.log
stderr_logfile=/tmp/nginx.err
